# udev/99-custom.rules
ACTION=="add", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c345", RUN+="/home/thecrawler/run_xmodmap.sh"
